# base/radix-system/etc/rc.d/rc.local
#!/bin/sh
#
# /etc/rc.d/rc.local  -  Local system initialization script.
#
# Put any local startup commands in here.
#

