# dev/perl/5.22.1-x86_32/profile.d/perl5-x32.sh.in
#!/bin/sh

PERL5LIB=/usr/lib@LIBDIRSUFFIX@/perl5:/usr/lib@LIBDIRSUFFIX@/perl5/vendor_perl:/usr/lib@LIBDIRSUFFIX@/perl5/site_perl
export PERL5LIB
